# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
directory = /tmp
childlogdir = /tmp

[program:traefik]
command=/home/user/traefik/traefik -c /home/user/traefik/traefik.toml
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:wsagent]
command=/home/user/tomcat8/bin/catalina.sh run
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:traefik_conf]
command=/home/user/traefik/traefik_conf.sh
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0



